# etc/services.d/miniboard/run
#!/usr/bin/with-contenv bash
exec miniboard
